# Bingo-Back/cloudflared-quick-tunnel.service
# Cloudflare quick tunnel - keeps backend reachable via HTTPS
# Copy to VPS: sudo cp cloudflared-quick-tunnel.service /etc/systemd/system/
# Then: sudo systemctl daemon-reload && sudo systemctl enable cloudflared-quick-tunnel && sudo systemctl start cloudflared-quick-tunnel
# Check: sudo systemctl status cloudflared-quick-tunnel
#
# Note: Quick tunnel URL changes on each restart. For a stable URL, use a named tunnel with a Cloudflare zone.

[Unit]
Description=Cloudflare quick tunnel to local backend
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel --url http://127.0.0.1:3001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
